# installers/v2/macos/pkg/scripts/postinstall
#!/bin/bash
# postinstall — runs after the Tadaima.app bundle is copied into
# /Applications by pkgbuild. Responsibilities (in order):
#
#   1. Select the correct Node runtime for the installing machine's arch
#      and symlink it as `runtime` inside the bundle.
#   2. Install @psychout98/tadaima from the bundled tarball into
#      Contents/Resources/agent/.
#   3. Run the bundled TadaimaConfig.app first-run GUI, waiting for it to
#      exit with code 0 before continuing.
#   4. Render com.tadaima.agent.plist.template into ~/Library/LaunchAgents
#      and bootstrap it with launchctl.
#
# Errors are logged to /tmp/tadaima-postinstall.log for diagnosis.
# Any non-zero exit aborts the installer.

set -euo pipefail

INSTALL_LOG="/tmp/tadaima-postinstall.log"
exec > >(tee -a "$INSTALL_LOG") 2>&1
echo ""
echo "=== tadaima postinstall started at $(date) ==="

APP="/Applications/Tadaima.app"
RES="$APP/Contents/Resources"

# ── Resolve the real GUI user ────────────────────────────────────────
# The postinstall runs as root. $USER is "root", not the logged-in
# user. Resolve the real GUI user from /dev/console ownership.
REAL_USER=$(/usr/bin/stat -f "%Su" /dev/console)
if [ "$REAL_USER" = "root" ] || [ -z "$REAL_USER" ]; then
    echo "[FATAL] could not determine GUI user (got '$REAL_USER')"
    exit 1
fi
USER_HOME=$(eval echo "~$REAL_USER")
USER_UID=$(id -u "$REAL_USER")
echo "[info] GUI user=$REAL_USER home=$USER_HOME uid=$USER_UID"

LOG_DIR="$USER_HOME/Library/Logs/Tadaima"
LOG_FILE="$LOG_DIR/agent.log"
PLIST_SRC="$RES/com.tadaima.agent.plist.template"
PLIST_DEST="$USER_HOME/Library/LaunchAgents/com.tadaima.agent.plist"

# Verify the app bundle landed where we expect.
if [ ! -d "$APP" ]; then
    echo "[FATAL] Tadaima.app not found at $APP"
    echo "        This means the installer placed the app somewhere else."
    echo "        Check that the distribution.xml uses enable_localSystem only."
    exit 1
fi

# ── 1. Select correct runtime for the machine's architecture ─────────
ARCH=$(uname -m)
echo "[info] architecture=$ARCH"
case "$ARCH" in
    arm64)  SRC_RUNTIME="runtime-arm64" ;;
    x86_64) SRC_RUNTIME="runtime-x64" ;;
    *)
        echo "[FATAL] unsupported architecture: $ARCH"
        exit 1
        ;;
esac

if [ ! -d "$RES/$SRC_RUNTIME" ]; then
    echo "[FATAL] runtime directory missing: $RES/$SRC_RUNTIME"
    echo "        Contents of $RES:"
    ls -la "$RES/" || true
    exit 1
fi

# Create the arch-neutral `runtime` symlink.
ln -sfn "$SRC_RUNTIME" "$RES/runtime"
echo "[info] runtime symlink: $RES/runtime -> $SRC_RUNTIME"

NODE="$RES/runtime/bin/node"
NPM_CLI="$RES/runtime/lib/node_modules/npm/bin/npm-cli.js"
if [ ! -x "$NODE" ]; then
    echo "[FATAL] bundled node missing or not executable: $NODE"
    exit 1
fi
if [ ! -f "$NPM_CLI" ]; then
    echo "[FATAL] bundled npm-cli.js missing: $NPM_CLI"
    exit 1
fi
echo "[info] node=$("$NODE" --version) npm=$("$NODE" "$NPM_CLI" --version)"

# ── 2. Install the pinned agent from the bundled tarball ─────────────
TARBALL="$RES/agent-tarball.tgz"
if [ ! -f "$TARBALL" ]; then
    echo "[FATAL] bundled agent tarball missing: $TARBALL"
    exit 1
fi

mkdir -p "$RES/agent"
echo "[info] running npm install from tarball (this may take a moment)..."
if ! "$NODE" "$NPM_CLI" install -g \
    --prefix "$RES/agent" \
    "$TARBALL" 2>&1; then
    echo ""
    echo "[FATAL] npm install failed. See output above."
    echo "        This usually means a network issue prevented npm from"
    echo "        downloading dependencies. Check your internet connection."
    exit 1
fi

AGENT_BIN="$RES/agent/bin/tadaima"
if [ ! -x "$AGENT_BIN" ]; then
    echo "[FATAL] agent binary missing after npm install: $AGENT_BIN"
    echo "        Contents of $RES/agent/bin/:"
    ls -la "$RES/agent/bin/" 2>/dev/null || echo "        (directory does not exist)"
    exit 1
fi
echo "[info] agent installed: $AGENT_BIN"

# ── 3. Launch the first-run config GUI ───────────────────────────────
sudo -u "$REAL_USER" mkdir -p "$LOG_DIR"

CONFIG_APP="$RES/TadaimaConfig.app"
if [ ! -d "$CONFIG_APP" ]; then
    echo "[FATAL] config GUI missing: $CONFIG_APP"
    exit 1
fi

# Remove macOS quarantine flags from unsigned apps so Gatekeeper doesn't
# block them when launched programmatically.
xattr -rd com.apple.quarantine "$APP" 2>/dev/null || true
xattr -rd com.apple.quarantine "$CONFIG_APP" 2>/dev/null || true

echo "[info] launching config GUI as $REAL_USER..."
# `open -W` blocks until the app terminates. launchctl asuser runs it in
# the installing user's GUI session so NSOpenPanel etc. work correctly.
if ! launchctl asuser "$USER_UID" sudo -u "$REAL_USER" open -W "$CONFIG_APP"; then
    echo "[FATAL] config GUI exited non-zero"
    echo "        The user may have cancelled, or the app crashed."
    echo "        Check Console.app for TadaimaConfig crash logs."
    exit 1
fi
echo "[info] config GUI completed successfully"

# ── 4. Install LaunchAgent plist ─────────────────────────────────────
sudo -u "$REAL_USER" mkdir -p "$(dirname "$PLIST_DEST")"
sed \
    -e "s#__AGENT_BIN__#$AGENT_BIN#g" \
    -e "s#__LOG_PATH__#$LOG_FILE#g" \
    "$PLIST_SRC" > "$PLIST_DEST"
chown "$REAL_USER" "$PLIST_DEST"
echo "[info] plist installed: $PLIST_DEST"

# ── 5. Bootstrap the LaunchAgent ─────────────────────────────────────
launchctl bootstrap "gui/$USER_UID" "$PLIST_DEST" 2>&1 || {
    # Non-fatal: the user can re-bootstrap from the menu bar Settings.
    echo "[warn] launchctl bootstrap failed (non-fatal, agent may already be loaded)"
}

echo ""
echo "=== tadaima postinstall completed at $(date) ==="
exit 0
